# misc/ID-Server/netcore-directory.service
[Unit]
Description=NetCore Directory Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/netcore-directory
ExecStart=/usr/bin/python3 /opt/netcore-directory/netcore_directory_server.py --host 0.0.0.0 --port 8095 --db /opt/netcore-directory/netcore_directory.db
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
