# system-backend/tts/netcore-piper.service
[Unit]
Description=NetCore local Piper TTS provider
After=network.target
Before=bluestation-bs.service

[Service]
Type=simple
User=bluestation
Group=bluestation
WorkingDirectory=/var/lib/netcore/piper
Environment=HOME=/var/lib/netcore/piper
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/netcore-piper/bin/python -m piper.http_server -m de_DE-thorsten-medium --data-dir /var/lib/netcore/piper --host 127.0.0.1 --port 5000
Restart=on-failure
RestartSec=3
TimeoutStopSec=15
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/netcore/piper /var/cache/netcore/tts

[Install]
WantedBy=multi-user.target
