# system-backend/mobility-core/systemd/netcore-mobility-core.service
[Unit]
Description=NetCore-Tetra Mobility Core (OPEN LAB MODE)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=netcore
Group=netcore
ExecStart=/usr/local/bin/netcore-mobility-core --config /etc/netcore/mobility-core.toml
Restart=on-failure
RestartSec=3
Environment=RUST_LOG=info
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/netcore-mobility-core

[Install]
WantedBy=multi-user.target
