# contrib/systemd/tetra.service.d/20-packet-data-gateway.conf
[Service]
# Required to create/configure the TUN interface and managed routes/firewall.
# Deliberately do not replace an existing CapabilityBoundingSet from the base unit.
AmbientCapabilities=CAP_NET_ADMIN
PrivateDevices=no
ProtectKernelTunables=no
RuntimeDirectory=netcore-tetra
RuntimeDirectoryMode=0750
ExecStopPost=-/usr/local/libexec/netcore-tetra-packet-gateway-cleanup
