# system-backend/media-switch/systemd/netcore-media-switch.service
[Unit]
Description=NetCore Media Switch (OPEN LAB)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=netcore
Group=netcore
ExecStart=/usr/local/bin/netcore-media-switch --config /etc/netcore/media-switch.toml
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/netcore-media-switch
Environment=RUST_LOG=netcore_media_switch=info

[Install]
WantedBy=multi-user.target
